# timetrack.service
[Unit]
Description=Timetracking daemon

[Service]
Type=simple
ExecStart=/opt/timetrack/src/timetrack.py

[Install]
WantedBy=multi-user.target
